# Frontend/Nginx/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log    	; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisor/               	; where child log files will live

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile = /var/log/supervisor/nginx-stderr.log
stdout_logfile = /var/log/supervisor/nginx-stdout.log

[program:web-socket]
command=/usr/sbin/node src/app.js
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile = /var/log/supervisor/web-socket-stderr.log
stdout_logfile = /var/log/supervisor/web-socket-stdout.log

[inet_http_server]
port = *:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

